# templates/default/tomcat.service.erb
# Systemd unit file for tomcat
[Unit]
Description=Apache Tomcat Web Application Container
After=syslog.target network.target

[Service]
Type=forking

Environment=JAVA_HOME=<%= node['jvogt_tomcat']['java_home'] %>
Environment=CATALINA_PID=<%= node['jvogt_tomcat']['install_path'] %>/temp/tomcat.pid
Environment=CATALINA_HOME=<%= node['jvogt_tomcat']['install_path'] %>
Environment=CATALINA_BASE=<%= node['jvogt_tomcat']['install_path'] %>
Environment='CATALINA_OPTS=-Xms512M -Xmx1024M -server -XX:+UseParallelGC'
Environment='JAVA_OPTS=-Djava.awt.headless=true -Djava.security.egd=file:/dev/./urandom'

ExecStart=<%= node['jvogt_tomcat']['install_path'] %>/bin/startup.sh
ExecStop=/bin/kill -15 $MAINPID

User=<%= node['jvogt_tomcat']['user'] %>
Group=<%= node['jvogt_tomcat']['group'] %>

[Install]
WantedBy=multi-user.target